# scripts/websocket-server@.service
[Unit]
Description=WebSocket Server for Sign Classification (Port %i)
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/team5-waterandfish-BE
Environment=PATH=/home/ubuntu/team5-waterandfish-BE/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/ubuntu/team5-waterandfish-BE/venv/bin/python3 /home/ubuntu/team5-waterandfish-BE/scripts/websocket_server_wrapper.py %i
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=websocket-server-%i

# 리소스 제한
LimitNOFILE=65536
TimeoutStartSec=60
TimeoutStopSec=30

# 프로세스 관리
KillMode=mixed
KillSignal=SIGTERM

# 보안 설정
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
